# systemd/nwall-lease-trigger.service
[Unit]
Description=nwall lease HTTP trigger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env /usr/local/bin/nwall lease trigger
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
