# systemd/nwall-lease.service
[Unit]
Description=nwall TCP lease agent
After=network-online.target nwall.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env /usr/local/bin/nwall lease agent
Restart=on-failure
RestartSec=2s
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
